# OS/RootFS/etc/init.d/boot
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org



	mknod /dev/gpio-key c 51 0
	mknod /dev/gpio-led c 50 0
	chmod 644 /dev/gpio-key
	chmod 644 /dev/gpio-led	




